# home/kruayd/.config/autostart/pipewire-pulse.desktop
[Desktop Entry]
Exec=pipewire-pulse
Icon=
Name=pipewire-pulse
Path=
Terminal=False
Type=Application
