# home/kruayd/.config/autostart/pipewire.desktop
[Desktop Entry]
Exec=pipewire
Icon=
Name=pipewire
Path=
Terminal=False
Type=Application
