# etc/sv/autorotate/run
#!/bin/sh
exec /usr/local/bin/autorotate.sh
